# covid-milchbartstrasse.ini
[program:covid-milchbartstrasse]
directory=%(ENV_HOME)s/apps/covid.milchbartstrasse.de
command=npm start
environment=NODE_ENV=production
